# files-SMP/etc/crontabs/root
@reboot /etc/red-led-off
@reboot /etc/init.d/adblock start &
*/30 * * * * /etc/speed